# systemd/dot-config/systemd/user/blueman.service
[Unit]
Description=Blueman Bluetooth Manager
PartOf=autostart.target

[Service]
ExecStart=blueman-applet
