# systemd/dot-config/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager Applet
PartOf=autostart.target

[Service]
ExecStart=nm-applet
